# systemd/user/break-reminder.service
[Unit]
Description=Break Reminder Notification

[Service]
Type=oneshot
ExecStart=/home/r4ppz/Arch-dotfiles/scripts/break-reminder.sh

